# cv-backend.service
[Unit]
Description=CV-Backend Service

[Service]
WorkingDirectory=<WorkDir>
ExecStart=<PythonDir> main.py
Restart=always